# conf/dot/.config/systemd/user/jn.service
[Unit]
Description=Jupyter Notebook

[Service]
ExecStart=/bin/bash "$HOME/bin/jn"
Restart=on-failure

[Install]
WantedBy=default.target
